# root/etc/udev/rules.d/81-wol.rules.tmpl
{{- $device := "eno1" -}}
ACTION=="add", SUBSYSTEM=="net", NAME=="{{ $device }}", RUN+="/usr/bin/ethtool -s {{ $device }} wol g"
